# gasmonitor.service
[Unit]
Description=Industrial Gas Monitoring Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/.venv/bin/python __INSTALL_DIR__/main.py --config __INSTALL_DIR__/config.ini
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
