# oxygen-monitor.service
[Unit]
Description=Raspberry Pi Oxygen Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/oxygen-monitor.conf
ExecStart=__INSTALL_DIR__/.venv/bin/python __INSTALL_DIR__/oxygen_monitor.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
